# root/usr/lib/systemd/system/paddleocr-sock.path
[Unit]
Description=PaddleOCR HTTP socket file

[Path]
PathChanged=/run/paddleocr

[Install]
WantedBy=multi-user.target
